# docker-assets/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:node]
autorestart=true
directory=/var/app/current
command=/bin/sh localdev.sh node .
